# deploy/launchd/com.kalshi-sentinel.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kalshi-sentinel.backend</string>

  <key>WorkingDirectory</key>
  <string>/Users/richard/.openclaw/workspace/kalshi-sentinel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd backend && source .venv/bin/activate && python app.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/richard/.openclaw/workspace/kalshi-sentinel/data/launchd-backend.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/richard/.openclaw/workspace/kalshi-sentinel/data/launchd-backend.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
